# templates/default/etc/profile.d/oozie.sh.erb

# this file is managed by chef

OOZIE_URL=http://<%= @oozie_node %>:<%= @oozie_port %>/oozie
export OOZIE_URL
