# templates/default/etc/profile.d/spark.sh.erb

PATH=$PATH:<%= @spark_path %>
export PATH
